# tmp/setbg.service
# Set backgrond
#

[Unit]
Description=Set backgrond

[Service]
ExecStart=/usr/bin/free

[Install]
WantedBy=default.target

