# config/systemd/user/clipmenud.service
[Unit]
Description=Clipmenu daemon
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecCondition=/usr/bin/bash -c '[[ "$XDG_SESSION_DESKTOP" =~ (qtile|i3|openbox) ]]'
ExecStart=/usr/bin/clipmenud
Restart=on-failure
RestartSec=5

Environment=DISPLAY=:0
# Environment=WAYLAND_DISPLAY=wayland-0

[Install]
WantedBy=graphical-session.target
