# config/systemd/user/clipmenud.service
#
# Service for Clipmenu
#

[Unit]
Description=Clipmenu daemon
After=default.target

[Service]
Type=simple
ExecCondition=/usr/bin/bash -c 'pgrep -x "qtile|i3|openbox" > /dev/null'
ExecStart=/usr/bin/clipmenud
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0
# Environment=WAYLAND_DISPLAY=wayland-0

[Install]
WantedBy=default.target

